# root/etc/s6-overlay/s6-rc.d/init-kimai-config/up
/etc/s6-overlay/s6-rc.d/init-kimai-config/run
